# deploy/joshlayani-api.service
[Unit]
Description=joshlayani.com analytics API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/joshlayani.com
Environment=PATH=/usr/bin:/bin
EnvironmentFile=/etc/joshlayani.com/api.env
ExecStart=/usr/bin/node /var/www/joshlayani.com/api/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
